# backend/scripts/com.cristopher.daily-signal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cristopher.daily-signal</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/cristovao/Documents/Projects/Cristopher/backend/scripts/run_daily_signal.sh</string>
    </array>

    <!-- Run every day at 09:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If Mac was asleep / off at 09:00, run as soon as it wakes (script self-skips if already ran today) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/cristovao/Documents/Projects/Cristopher/backend/scripts/daily_signal.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/cristovao/Documents/Projects/Cristopher/backend/scripts/daily_signal.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/cristovao/Documents/Projects/Cristopher/backend</string>
</dict>
</plist>
